# tts.service
[Unit]
Description=Stream TTS

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=bash start.sh
StandardOutput=null
StandardError=null
WorkingDirectory=/home/alexandra/git/streamtts

[Install]
WantedBy=default.target
